# root-image/etc/rc.local
# powersave
powerdown >& /var/log/powerdown

# login
exec login -f march
